# xdoc/check_service.cron
*/1 * * * * ./check_service.sh ztf 8848 >> check_service.log